# src/unizab.etc.cron.d
#
# $Header$
#
0 */1  * * *	root [ -x /var/opt/UNIzab/bin/col2za ] && /var/opt/UNIzab/bin/col2za hour
59 23  * * *	root [ -x /var/opt/UNIzab/bin/col2za ] && /var/opt/UNIzab/bin/col2za day
